# daemon/rc-script/xapsd
#!/bin/sh
# $FreeBSD
#

# PROVIDE: xapsd
# REQUIRE: DAEMON  @comment  %%CLAMD%%
# BEFORE: mail
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf[.local] to enable xapsd:
# xapsd_enable="YES"
#

. /etc/rc.subr

name="xapsd"
rcvar=xapsd_enable

load_rc_config $name

#defaults
: ${xapsd_enable="NO"}
: ${xapsd_pidfile:-/var/run/dovecot/xapsd.pid}
: ${xapsd_configname="xapsd.yaml"}
: ${xapsd_configpath="/usr/local/etc/xapsd"}

command=/usr/local/bin/${name}

#required_dirs=/var/db/xapsd
required_files=${xapsd_configpath}/${xapsd_configname}

command_args="-configPath ${xapsd_configpath} -configName ${xapsd_configname}"

run_rc_command "$1"
